# deploy/supervisor/quickshare-worker.conf
; QuickShare Laravel Queue Worker
; Place this file in /etc/supervisor/conf.d/ (or /etc/supervisord.d/)
; Adjust user/paths for your server.

[program:quickshare-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/quickshare/artisan queue:work --queue=notifications,default --sleep=3 --tries=3 --max-jobs=1000 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/quickshare/storage/logs/worker.log
stopwaitsecs=3600
environment=HOME="/var/www"

; Optional: load .env directly if Supervisor cannot see the shell env.
; environment=HOME="/var/www",APP_ENV="production",APP_DEBUG="false"
